# deploy/config/crawler.timer
[Unit]
Description=Run crawler every 30 min

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target

